# app-emulation/hasp/files/hasplm.init-8.30
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend () {
    need net localmount
    after aksusbd
}

start() {
	ebegin "Starting HASP License manager"
	start-stop-daemon --start --quiet  --exec /usr/sbin/hasplm -- ${HASPLM_ARGS}
	eend $?
}

stop() {
	ebegin "Stopping HASP License manager"
	start-stop-daemon --stop --quiet --retry 10 --exec /usr/sbin/hasplm
	eend $?
}
